# warehouse.service
[Unit]
Description=Warehouse Management System
Documentation=https://github.com/yourusername/warehouse
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ruslan
Group=ruslan
WorkingDirectory=/home/ruslan/Desktop/Main-project-13strategyexperts
Environment="NODE_ENV=production"
Environment="PATH=/home/ruslan/.nvm/versions/node/v20.18.0/bin:/home/ruslan/.nvm/versions/node/v20.17.0/bin:/home/ruslan/.nvm/versions/node/v20.16.0/bin:/usr/bin:/usr/local/bin:/bin"

# Загрузка переменных из .env
EnvironmentFile=/home/ruslan/Desktop/Main-project-13strategyexperts/.env

# Команда запуска (npm и node будут найдены автоматически через PATH)
ExecStart=/usr/bin/env npm run dev

# Логирование
StandardOutput=append:/var/log/warehouse/stdout.log
StandardError=append:/var/log/warehouse/stderr.log
SyslogIdentifier=warehouse

# Перезапуск при сбое
Restart=always
RestartSec=10

# Ограничения безопасности
NoNewPrivileges=true
PrivateTmp=true

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
